# src/scripts/etc/supervisor/conf.d/snf-cdmi.conf
[program:snf-cdmi]
command=/root/snf-cdmi/start-snf-cdmi.supervisord

